# infra/gpu-scheduler/gpu-scheduler.service
[Unit]
Description=GPU Priority Scheduler
Documentation=https://github.com/yuribarreira/apps-to-make-money
After=docker.service ollama.service
Requires=docker.service

[Service]
Type=simple
User=yuri
WorkingDirectory=/home/yuri/apps-to-make-money/infra/gpu-scheduler
EnvironmentFile=/home/yuri/income-services/shared/.gpu-scheduler.env
ExecStart=/usr/bin/python3 /home/yuri/apps-to-make-money/infra/gpu-scheduler/scheduler.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
